# scripts/crontab
*/1 * * * * python /usr/src/app/heyGiulia.py
